# config/etc/periodic/daily/clean_tmp
#!/bin/sh

find /tmp/* -mtime +2 -exec rm {} \;
